# tpm2keyunlock.service
[Unit]
Description=TPM2 post-install manager

[Service]
User=root
WorkingDirectory=LOCATION
ReadWriteDirectories=LOCATION
TimeoutSec=180
ExecStart=LOCATION/tpm2PolicyConfig

[Install]
WantedBy=multi-user.target
